# src/infrastructure/eventops_gateway.service
[Unit]
Description=Sovereign EventOps Rust/Python Validation Gateway
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/code/src/rust/eventops_pyo3
# Run the gateway strictly through the virtual environment containing the maturin build
ExecStart=/home/ubuntu/code/src/rust/eventops_pyo3/.venv/bin/python gateway_server.py
Restart=on-failure
RestartSec=5s

# Security boundaries
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
