# scripts/monitoring/dashboard.service
# Systemd service for INFRA-2 Monitoring Dashboard
# Deploy to: /etc/systemd/system/agentic-dashboard.service
# Commands:
#   sudo systemctl enable agentic-dashboard
#   sudo systemctl start agentic-dashboard
#   sudo systemctl status agentic-dashboard

[Unit]
Description=Agentic Flow Real-Time Monitoring Dashboard
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/agentic-flow
ExecStart=/usr/bin/node /home/ubuntu/agentic-flow/scripts/monitoring/dashboard_server.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=DASHBOARD_PORT=8080

# Logging
StandardOutput=append:/home/ubuntu/dashboard.log
StandardError=append:/home/ubuntu/dashboard.log

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu

[Install]
WantedBy=multi-user.target

