# investing/agentic-flow/scripts/infra/com.agentic-flow.drift-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentic-flow.drift-check</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/shahroozbhopti/Documents/code/investing/agentic-flow/scripts/infra/drift-check.sh</string>
        <string>--log</string>
        <string>/Users/shahroozbhopti/Documents/code/investing/agentic-flow/scripts/infra/logs/drift-check.log</string>
    </array>

    <!-- Run every 4 hours (14400 seconds) -->
    <key>StartInterval</key>
    <integer>14400</integer>

    <key>StandardOutPath</key>
    <string>/Users/shahroozbhopti/Documents/code/investing/agentic-flow/scripts/infra/logs/drift-check-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shahroozbhopti/Documents/code/investing/agentic-flow/scripts/infra/logs/drift-check-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
