# _SYSTEM/_AUTOMATION/com.bhopti.legal.wsjf-escalator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bhopti.legal.wsjf-escalator</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/shahroozbhopti/Documents/code/investing/agentic-flow/scripts/validators/semantic/validate-wsjf-escalation.sh</string>
    </array>
    
    <key>StartInterval</key>
    <integer>900</integer>
    
    <!-- Alternative: Run every 15 minutes (900 seconds) -->
    <!-- To switch back to daily at 9 AM, use StartCalendarInterval instead -->
    
    <key>StandardOutPath</key>
    <string>/Users/shahroozbhopti/Library/Logs/wsjf-escalator.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/shahroozbhopti/Library/Logs/wsjf-escalator-error.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
