# scripts/deploy/systemd_dashboard.service
[Unit]
Description=Agentic Flow Web Dashboard
After=network.target

[Service]
Type=notify
User=agentic-flow
Group=agentic-flow
WorkingDirectory=/opt/agentic-flow
Environment="PROJECT_ROOT=/opt/agentic-flow"
Environment="PYTHONPATH=/opt/agentic-flow/scripts"
Environment="FLASK_SECRET_KEY=<REPLACE_WITH_SECURE_KEY>"
ExecStart=/opt/agentic-flow/venv/bin/gunicorn \
    --bind 127.0.0.1:8000 \
    --workers 4 \
    --worker-class eventlet \
    --timeout 120 \
    --access-logfile /var/log/agentic-flow/access.log \
    --error-logfile /var/log/agentic-flow/error.log \
    --log-level info \
    'scripts.web_dashboard:app'
    
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
KillSignal=SIGQUIT
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
