# scripts/launchagents/com.bhopti.roam.staleness.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bhopti.roam.staleness.watchdog</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/shahroozbhopti/Documents/code/investing/agentic-flow/scripts/validators/roam-staleness-watchdog.sh</string>
    </array>
    
    <key>StartInterval</key>
    <integer>3600</integer>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>/Users/shahroozbhopti/Library/Logs/roam-staleness-watchdog.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/shahroozbhopti/Library/Logs/roam-staleness-watchdog-error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/shahroozbhopti/Documents/code/investing/agentic-flow</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>LowPriorityIO</key>
    <true/>
    
    <key>Nice</key>
    <integer>1</integer>
</dict>
</plist>
