# projects/agentic-flow-core/scripts/infra/cpanel/com.agentic-flow.ssl-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentic-flow.ssl-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>source $HOME/Documents/code/agentic-flow-core/config/.env.cpanel 2>/dev/null; $HOME/Documents/code/agentic-flow-core/scripts/infra/cpanel/ssl-monitor.sh once</string>
    </array>

    <!-- Run every 4 hours -->
    <key>StartInterval</key>
    <integer>14400</integer>

    <key>StandardOutPath</key>
    <string>/tmp/ssl-monitor.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/ssl-monitor.err</string>

    <key>RunAtLoad</key>
    <true/>

    <!-- Don't retry on failure -->
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
