# config/deployment/swarm-api-server.service
[Unit]
Description=Sovereign Swarm SSR API Server (Bridge)
After=network.target

[Service]
Type=simple
User=ubuntu
# Replace /var/www/agentic-flow with the absolute path to the repository on target host
WorkingDirectory=/var/www/agentic-flow/swarm-core-app
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=swarm-api-server
Environment=NODE_ENV=production
# Load environment secrets if available
EnvironmentFile=-/var/www/agentic-flow/.env
EnvironmentFile=-/var/www/agentic-flow/.env.integration

[Install]
WantedBy=multi-user.target
