# deploy_package/configs/trading-api.service
[Unit]
Description=Agentic Flow Trading API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/agentic-flow/trading
Environment="NODE_ENV=production"
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
