# config/monitoring/gitlab_health_check.cron
# GitLab Health Check - Run every minute
* * * * * /Users/shahroozbhopti/Documents/code/investing/agentic-flow/config/monitoring/gitlab_health_check.sh >> /var/log/gitlab-health.log 2>&1

# Daily backup verification - Run at 3 AM
0 3 * * * ssh gitlab.interface.splitcite.com 'ls -la /var/opt/gitlab/backups/*.tar 2>/dev/null | tail -1' >> /var/log/gitlab-backup-check.log 2>&1
